# deploy/legal-wecom-backup.timer
[Unit]
Description=Daily Legal WeCom Automation backup

[Timer]
OnCalendar=*-*-* 02:30:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
