# deploy/legal-wecom-backup.service
[Unit]
Description=Legal WeCom Automation consistent backup
After=local-fs.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/legal-wecom-automation
EnvironmentFile=/opt/legal-wecom-automation/.env
ExecStart=/opt/legal-wecom-automation/.venv/bin/python /opt/legal-wecom-automation/scripts/backup.py
